# matrix-presence-bridge.service
[Unit]
Description=Matrix Presence Bridge (System-Idle → Synapse)
After=network-online.target graphical-session.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/matrix-presence-bridge.sh
Restart=always
RestartSec=30
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority

[Install]
WantedBy=default.target
